# examples/n8n.service
[Unit]
Description=n8n workflow automation tool
After=network.target

[Service]
Type=simple
User=n8n
Group=n8n
WorkingDirectory=/home/n8n
ExecStart=/usr/bin/npx n8n start
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=WEBHOOK_URL=https://yourdomain.com/
Environment=GENERIC_TIMEZONE=Europe/London
Environment=N8N_PORT=5678
Environment=N8N_PROTOCOL=http
Environment=N8N_HOST=localhost

# Optional: Enable basic authentication
# Environment=N8N_BASIC_AUTH_ACTIVE=true
# Environment=N8N_BASIC_AUTH_USER=admin
# Environment=N8N_BASIC_AUTH_PASSWORD=secure_password_here

# Optional: Database configuration (PostgreSQL)
# Environment=DB_TYPE=postgresdb
# Environment=DB_POSTGRESDB_HOST=localhost
# Environment=DB_POSTGRESDB_DATABASE=n8n
# Environment=DB_POSTGRESDB_USER=n8n_user
# Environment=DB_POSTGRESDB_PASSWORD=your_password

# Optional: Email configuration
# Environment=N8N_EMAIL_MODE=smtp
# Environment=N8N_SMTP_HOST=smtp.gmail.com
# Environment=N8N_SMTP_PORT=587
# Environment=N8N_SMTP_USER=your-email@gmail.com
# Environment=N8N_SMTP_PASS=your-app-password

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/n8n

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
